# util/server/crontabs/goaccess
# /etc/cron.d/

### run goaccess command every 24 hours
0 */24 * * * root cd /vol/drafty_source/drafty/util/server/ && ./run_goaccess
